# systemd/aria-gateway.service
[Unit]
Description=ARIA Gateway Daemon (Telegram)
After=default.target
Wants=default.target

[Service]
Type=notify
WorkingDirectory=%h/coding/aria
EnvironmentFile=%h/coding/aria/.env
Environment="ARIA_HOME=%h/coding/aria"
Environment="ARIA_RUNTIME=%h/coding/aria/.aria/runtime"
Environment="KILOCODE_CONFIG_DIR=%h/coding/aria/.aria/kilocode"
ExecStart=%h/coding/aria/.venv/bin/python -m aria.gateway.daemon
Restart=on-failure
RestartSec=5s
WatchdogSec=60s
NotifyAccess=main

# Security hardening (desktop-compatible baseline)
# NOTE 1: In systemd --user sessions on desktop Linux, both PrivateDevices=true
# and ProtectKernelModules=true fail at exec with status=218/CAPABILITIES
# ("Failed to drop capabilities: Operation not permitted"). Keep those
# directives disabled for user mode; document divergence in ADR-0008.
# NOTE 2: Gateway spawns KiloCode/Node child processes through ConductorBridge.
# Node/V8 requires executable memory pages; MemoryDenyWriteExecute=true causes
# V8 fatal crashes during child session spawn. Keep it disabled for gateway.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=tmpfs
ReadWritePaths=%h/coding/aria/.aria/runtime %h/coding/aria/.aria/kilocode/sessions
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native

[Install]
WantedBy=default.target
