# systemd/aria-gateway.service
[Unit]
Description=ARIA Gateway Daemon (Telegram)
After=default.target
Wants=default.target

[Service]
Type=notify
WorkingDirectory=%h/coding/aria
EnvironmentFile=%h/coding/aria/.env
Environment="ARIA_HOME=%h/coding/aria"
Environment="ARIA_RUNTIME=%h/coding/aria/.aria/runtime"
Environment="KILOCODE_CONFIG_DIR=%h/coding/aria/.aria/kilocode"
ExecStart=%h/coding/aria/.venv/bin/python -m aria.gateway.daemon
Restart=on-failure
RestartSec=5s
WatchdogSec=60s
NotifyAccess=main

# Security hardening (user service baseline - minimal)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=tmpfs
ReadWritePaths=%h/coding/aria/.aria
PrivateTmp=true

[Install]
WantedBy=default.target
